# deploy/systemd/gx-opendtu-zero-export.service
[Unit]
Description=Zero-injection PV controller (OpenDTU + Cerbo GX via Modbus TCP) - Java
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gx-opendtu
WorkingDirectory=/opt/gx-opendtu
ExecStart=/usr/bin/java -jar /opt/gx-opendtu/gx-opendtu-java.jar --config /etc/gx-opendtu/config.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
